# resources/lumen.service
[Unit]
Description=Lumen Desktop Shell
PartOf=graphical-session.target
After=graphical-session.target
StartLimitIntervalSec=30
StartLimitBurst=5

[Service]
ExecStart=/usr/bin/lumen shell
Restart=on-failure
RestartSec=3
Slice=session.slice

[Install]
WantedBy=graphical-session.target
